# script/old/dbus.socket
[Unit]
Description=D-Bus Message Bus Socket
Before=sockets.target

[Socket]
ListenStream=%t/dbus_user_socket

[Install]
WantedBy=sockets.target
