# syst/config/files/etc/udev/rules.d/52-suspend.rules
#auto suspend
SUBSYSTEM=="power_supply", ATTR{status}=="Discharging", ATTR{capacity}=="[012345]", RUN+="/usr/bin/systemctl suspend"
