# config/systemd/user/basic.target
[Unit]
Description=Services needed on console
Wants=init.target
AllowIsolate=yes

[Install]
Alias=default.target
