# config/systemd/user/crontab.target
[Unit]
Description=Timers
AllowIsolate=yes
