# config/systemd/user/init.target
[Unit]
Description=Init user session
Wants=sockets.target
AllowIsolate=yes

[Install]
Alias=default.target
