# config/systemd/user/internet.target
[Unit]
Description=Services needing internet
AllowIsolate=yes
Wants=unlock.target
