# config/systemd/user/ssh-agent.socket
[Unit]
Description=ssh-agent socket

[Socket]
ListenStream=/run/user/%U/ssh_auth_sock

[Install]
WantedBy=sockets.target
