# syst/generate/files/0ADurin/systemd/system/user@.service.d/env.conf
[Service]
Environment=XDG_RUNTIME_DIR=/run/user/%I
