# config/.config/systemd/user/crontab-backup.timer
[Unit]
Description=Git backup Timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h

[Install]
WantedBy=crontab.target
