# config/.config/systemd/user/crontab-backup_push.timer
[Unit]
Description=Git backup push Timer

[Timer]
OnBootSec=30min
OnUnitActiveSec=3h

[Install]
WantedBy=crontab.target
