# config/.config/systemd/user/crontab-backup.service
[Unit]
Description=Git backup

[Service]
Type=simple
ExecStart=/bin/sh -c 'exec %h/mine/script/gitbackup -v save'
