# config/.config/systemd/user/init.target
[Unit]
Description=Init user session
Wants=init.target
AllowIsolate=yes

[Install]
Alias=default.target
