# config/.config/systemd/user/session.target
[Unit]
Description=user session
AllowIsolate=yes
Wants=basic.target
Wants=crontab.target

[Install]
Alias=default.target
