# config/.config/systemd/user/xsession@.target
[Unit]
Description=xsession on display %I
Wants=session.target
AllowIsolate=yes

[Install]
Alias=default.target
