# user/config/systemd/user/core.target
[Unit]
Description=Core Session
Wants=init.target
After=init.target
AllowIsolate=yes

[Install]
Alias=default.target
