# user/config/systemd/user/init.target
[Unit]
Description=Initialize Session
Wants=basic.target
AllowIsolate=yes

[Install]
Alias=default.target
