# user/config/systemd/user/network.target
[Unit]
Description=Internet Services
AllowIsolate=yes
