# user/config/systemd/user/unlock.target
[Unit]
Description=Unlock
AllowIsolate=yes
