# user/crontab/@boot-phare
#!/usr/bin/zsh

exec $HOME/mine/user/crontab/@boot --init phare
