# user/crontab/@boot-imb
#!/usr/bin/zsh

exec $HOME/mine/user/crontab/@boot imb
