# user/crontab/imb-unlock.sh
#!/usr/bin/zsh

exec $HOME/mine/user/crontab/@boot --unlock imb-unlock
